# deploy/systemd/dbeaver-mistral-proxy.service
[Unit]
Description=DBeaver -> Mistral proxy (OpenAI Responses compatible)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cloud
Group=cloud
WorkingDirectory=/home/cloud/not-safe/github/dbeaver-proxy-to-mistral
EnvironmentFile=/home/cloud/not-safe/github/dbeaver-proxy-to-mistral/.env
ExecStart=/home/cloud/not-safe/github/dbeaver-proxy-to-mistral/.venv/bin/python -m dbeaver_mistral_proxy
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
